# systemd/file-validator.service
[Unit]
Description=File Magic Number Validator
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /usr/local/bin/file-validator.py
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
